# ansible/templates/aria2.socket.j2
[Unit]
Description=aria2 socket activation

[Socket]
ListenStream=127.0.0.1:6800

[Install]
WantedBy=sockets.target
