# ansible/templates/gost.service.j2
[Unit]
Description=gost
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/gost {{ gost_arguments }}
Restart=always

[Install]
WantedBy=multi-user.target
