# home/dot_config/cron/crontab
33 * * * * fcrontab - <~/.config/cron/crontab
*/5 * * * * ~/.bin/ddns.sh
43 * * * * ~/.bin/artpip.sh new
